# deployments/single-vps/systemd/gateway.service
[Unit]
Description=Realtime Room Middleware - Gateway
After=network.target

[Service]
Type=simple
User=realtime
Group=realtime
WorkingDirectory=/opt/realtime-server/current
ExecStart=/opt/realtime-server/current/gateway
Environment=CONFIG_PATH=/opt/realtime-server/config/production.yaml
Restart=always
RestartSec=3
LimitNOFILE=1048576

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/realtime-server/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gateway

[Install]
WantedBy=multi-user.target
